# .config/crontab
0 0 * * * [ -f ~/.python_history] && rm ~/.python_history
